# deploy/semantic-index-build.service
# Out-of-process nightly graph rebuild conductor (WXYC/semantic-index#347).
# Oneshot service invoked by semantic-index-build.timer. Drives the round-trip:
# snapshot prod DB -> S3 seed -> Fargate build -> validate -> atomic swap.
#
# Install on the BS EC2 host:
#   sudo cp deploy/semantic-index-build.service /etc/systemd/system/
#   sudo cp deploy/semantic-index-build.timer   /etc/systemd/system/
#   sudo cp scripts/ec2-build-conductor.sh /home/ec2-user/
#   # edit the Environment= lines below with the stack-output SUBNETS / BUILD_SG
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now semantic-index-build.timer
#
# One-off manual run:  sudo systemctl start semantic-index-build.service
# Logs:                journalctl -u semantic-index-build.service -f
[Unit]
Description=semantic-index out-of-process nightly graph rebuild
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=ec2-user
# Fill these from the CloudFormation stack outputs (infra/build-job.yaml):
Environment=SUBNETS=subnet-xxxxxxxx
Environment=BUILD_SG=sg-xxxxxxxx
Environment=AWS_REGION=us-east-1
# Optional overrides (defaults shown in scripts/ec2-build-conductor.sh):
# Environment=DATA_DIR=/home/ec2-user/semantic-index-data
# Environment=IMAGE=semantic-index:latest
# Environment=MIN_ARTISTS=1000
ExecStart=/home/ec2-user/ec2-build-conductor.sh
# Comfortably above the conductor's BUILD_WAIT_CEILING_SECS (2700s) plus the
# snapshot/upload/download/validate I/O on the ~1.5 GB DB.
TimeoutStartSec=3600
